# conf/supervisor.ini
[program:gunicorn]
command=/home/chi/.pyenv/shims/pipenv run gunicorn -c gunicorn.py manage:app
directory=/home/chi/MrChiBlog
autostart=true
autorestart=true
startsecs=5
startretries=3
stopasgroup=true
killasgroup=true
user=chi
stdout_logfile=/mrchiblog/logs/gunicorn_access.log
stderr_logfile=/mrchiblog/logs/gunicorn_error.log


[program:rq-update]
command=/home/chi/.pyenv/shims/pipenv run flask rq worker update
directory=/home/chi/MrChiBlog
autostart=true
autorestart=true
startsecs=5
startretries=3
stopasgroup=true
killasgroup=true
user=chi
redirect_stderr=true
stdout_logfile=/mrchiblog/logs/rq_update.log


[program:rq-notify]
command=/home/chi/.pyenv/shims/pipenv run flask rq worker notify
directory=/home/chi/MrChiBlog
autostart=true
autorestart=true
startsecs=5
startretries=3
stopasgroup=true
killasgroup=true
user=chi
redirect_stderr=true
stdout_logfile=/mrchiblog/logs/rq_notify.log


[group:mrchiblog]
programs=gunicorn, rq-update, rq-notify
